# scripts/packages/nginx-agent.service
[Unit]
Description=NGINX Agent
Documentation=https://www.nginx.com/products/nginx-agent/
After=network.target
Wants=network-online.target
StartLimitIntervalSec=15
StartLimitBurst=5

[Service]
User=${AGENT_USER}
Group=${AGENT_GROUP}
Type=simple

ExecStartPre=/bin/mkdir -p ${AGENT_RUN_DIR}
ExecStartPre=/bin/mkdir -p ${AGENT_LOG_DIR}

ExecStart=${AGENT_EXE}

RestartSec=5
Restart=always

PermissionsStartOnly=true
PIDFile=${AGENT_RUN_DIR}/nginx-agent.pid
Environment=

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=nginx-agent

ExecStop=/bin/kill -2 $MAINPID
ExecStop=/bin/sleep 3

[Install]
WantedBy=multi-user.target
