# scripts/packages/nginx-agent.openrc
#!/sbin/openrc-run

description="NGINX Agent"
command_background=true

cfgfile=${cfgfile:-/etc/nginx-agent/nginx-agent.conf}
pidfile=/var/run/nginx-agent.pid
command=/usr/bin/nginx-agent
command_args=""
required_files="$cfgfile"

depend() {
    need net
    use dns logger netmount
}
