# scripts/packages/nginx-agent.service
[Unit]
Description=NGINX Agent
Documentation=https://github.com/nginx/agent#readme
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=15
StartLimitBurst=5

[Service]
User=${AGENT_USER}
Group=${AGENT_GROUP}
Type=simple

ExecStartPre=/bin/mkdir -p ${AGENT_RUN_DIR} ${AGENT_LOG_DIR}
ExecStartPre=/bin/chown -R ${AGENT_USER}:${AGENT_GROUP} ${AGENT_RUN_DIR} ${AGENT_LOG_DIR}
ExecStart=${AGENT_EXE}

RestartSec=5
Restart=on-failure
PermissionsStartOnly=true

KillSignal=SIGINT
TimeoutStopSec=5s

NoNewPrivileges=true
ProtectSystem=yes
ProtectHome=off
PrivateTmp=true

[Install]
WantedBy=multi-user.target
