# scripts/packages/nginx-agent.service
[Unit]
Description=NGINX Agent
Documentation=https://github.com/nginx/agent#readme
After=network.target
Wants=network.target
StartLimitIntervalSec=15
StartLimitBurst=5

[Service]
User=${AGENT_USER}
Group=${AGENT_GROUP}
Type=simple

LogsDirectory=nginx-agent
RuntimeDirectory=nginx-agent
StateDirectory=nginx-agent

ExecStart=${AGENT_EXE}
RestartSec=5
Restart=on-failure
PermissionsStartOnly=true
TimeoutStopSec=5s

NoNewPrivileges=true
ProtectHome=off
PrivateTmp=true

[Install]
WantedBy=multi-user.target
